# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/femto-usb-config-tool.service
[Unit]
Description=Foxbuntu USB Configuration Tool
After=femto-runonce.service meshtasticd.service
Before=serial-getty@ttyFIQ0.service
DefaultDependencies=no

[Service]
ExecStart=/usr/local/bin/femto-usb-config-tool.sh
Type=oneshot
RemainAfterExit=true
StandardOutput=journal+console
StandardError=journal+console
TTYPath=/dev/console

[Install]
WantedBy=serial-getty@ttyFIQ0.service

