# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/femto-boot-complete.service
[Unit]
Description=Boot Complete
DefaultDependencies=no
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/femto-boot-complete.sh
Type=oneshot
Restart=no
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=default.target
