# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/femto-watchclock-dog.service
[Unit]
Description=Monitor time jumps and restart meshtasticd
After=network.target

[Service]
ExecStart=/usr/local/bin/femto-watchclock-dog.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target

