# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/femto-usb-config-tool.service
[Unit]
Description=Foxbuntu USB Configuration Tool
After=meshtasticd.service
Before=serial-getty@ttyFIQ0.service
DefaultDependencies=no

[Service]
ExecStart=/usr/local/bin/femto-usb-config-tool.sh
Type=oneshot
Restart=no
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=default.target

