# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/femto-wifi-mesh-control.service
[Unit]
Description=Meshtastic control of wifi
After=meshtasticd.service

[Service]
ExecStart=/usr/local/bin/femto-wifi-mesh-control.sh
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
