# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/time-monitor.service
[Unit]
Description=Monitor time changes and restart meshtasticd
After=network.target

[Service]
ExecStart=/usr/local/bin/time_monitor.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target

