# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/ttyd.service
[Unit]
Description=ttyd Web Terminal
After=syslog.target
After=network.target

[Service]
ExecStart=/opt/ttyd/ttyd -W -6 --ssl --ssl-key /etc/ssl/private/ttyd.key --ssl-cert /etc/ssl/certs/ttyd.crt -i 0.0.0.0 login
Type=simple
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target
