# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/lib/systemd/system/wifi-mesh-control.service
[Unit]
Description=Meshtastic control of wifi
After=meshtasticd.service
Requires=meshtasticd.service

[Service]
ExecStart=/usr/local/bin/wifi-mesh-control.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
